# monitoring/health_reporting_scheduler.service
[Unit]
Description=CBS Health Reporting Scheduler Service
After=network.target

[Service]
User=cbs_user
Group=cbs_user
WorkingDirectory=/path/to/CBS_PYTHON
ExecStart=/usr/bin/python3 /path/to/CBS_PYTHON/monitoring/health_reporting_scheduler.py
Restart=on-failure
Environment=CBS_ADMIN_BASE_URL=http://localhost:8000/api/admin
Environment=CBS_ADMIN_API_KEY=your_api_key_here
Environment=PYTHONPATH=/path/to/CBS_PYTHON

[Install]
WantedBy=multi-user.target
